# laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/vhosts/alunos.profissionalizaead.com.br/httpdocs/artisan queue:listen --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=vagrant
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/vhosts/alunos.profissionalizaead.com.br/httpdocs/storage/logs/worker.log
stopwaitsecs=30